# nat-reflection/nat-reflection.service
[Unit]
Description=NAT reflection service

[Service]
ExecStart=/usr/local/bin/nat-reflection reload

[Install]
WantedBy=multi-user.target